# deploy/systemd/babylink-napcat-watchdog.service
[Unit]
Description=BabyLink NapCat online watchdog
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh /opt/babylink/deploy/napcat-watchdog.sh